# supervisord.conf
# supervisor + uwsgi + docker
[unix_http_server]
file=/var/run/supervisor.sock

[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[unix_http_server]
file=/var/run/supervisor.sock

;[program:uwsgi]
;autorestart=True
;autostart=True
;redirect_stderr=True
;startsecs=10
;satrtretries=3
;command=uwsgi --ini /shop_admin/uwsgi.ini
;directory=/shop_admin
;stdout_logfile_maxbytes = 50MB
;stdout_logfile_backups = 10
;stdout_logfile = /shop_admin/log/uwsgi_stdout.log
;stopasgroup=true
;killasgroup=true
;processes = 4
;max-requests=5000
# daemonize=/var/log/uwsgi/shop_admin.log
;chmod-socket=666
;vacuum=True


[program:daphne]
autorestart=True
autostart=True
redirect_stderr=True
startsecs=10
satrtretries=3
command=daphne -b 0.0.0.0 -p 8001 /shop_admin/adminDemo.asgi:application
directory=/shop_admin
stdout_logfile = /shop_admin/log/uwsgi_stdout.log
processes = 4
max-requests=5000
